# IP Monitor/ip_monitor
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

START=99
STOP=99

# Configuration for the executable script
SERVICE_PROG="/usr/bin/ip_monitor.sh"
SERVICE_NAME="ip_monitor"
INTERVAL=10 # Execution interval in seconds

start() {
    # It's good practice to ensure the executable is there before starting the loop
    if [ ! -x "$SERVICE_PROG" ]; then
        echo "Error: $SERVICE_PROG not found or not executable."
        exit 1
    fi
    
    # procd_open_instance is not strictly necessary for the looping style, 
    # but we include it for structure and logging setup.
    procd_open_instance "$SERVICE_NAME"
    procd_set_param command /bin/ash "$SERVICE_PROG"
    procd_set_param stdout 1 # Log output to system log
    procd_set_param stderr 1 # Log errors to system log
    procd_set_param respawn 5 5 5 # Respawn settings (though the loop handles restarts)
    procd_close_instance

    # Main loop that keeps the script running every 10 seconds.
    # The '&' runs the loop in the background so the init process can complete.
    while :; do
        "$SERVICE_PROG"
        sleep "$INTERVAL"
    done &
}

stop() {
    # Use the executable filename to ensure the correct background process is killed.
    killall "$SERVICE_NAME"
}